# daemon/ImageServer.service
[Unit]
Description=Image Server

[Service]
Type=forking
ExecStart=/usr/bin/ImageServer
Restart=on-failure

[Install]
WantedBy=multi-user.target